# src/templates/project_dir/examples/systemd/etc/systemd/system/[% project.package_name %].timer
[Unit]
Description=[% project.summary %]

[Timer]
Unit=[% project.package_name %].service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
